# supervisord.conf
[supervisord]
nodaemon=true

[program:redis-server]
command=/usr/bin/redis-server /etc/redis/redis.conf
pidfile=/var/run/redis/redis-server.pid
startsecs=5
autostart=true
autorestart=true


[program:mysql]
command=/usr/bin/pidproxy /var/run/mysqld/mysqld.pid /usr/bin/mysqld_safe --bind-address=0.0.0.0
startsecs=5
autostart=true
autorestart=true

[program:php-fpm]
command=/usr/sbin/php-fpm7.0 --nodaemonize
autostart=true
autorestart=true

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true


